# profile.d/gamescope.sh
#!/bin/bash

#sudo bash << EOF
#echo "#!/bin/bash" > /etc/profile.d/gamescope.sh
#echo "gamescope -e -- steam -steamos -pipewire-dmabuf -gamepadui" >> /etc/profile.d/gamescope.sh
#EOF
